# scripts/udev-rules/sigil-daemon.service
[Unit]
Description=Sigil MPC Signing Daemon
Documentation=https://github.com/chippr-robotics/sigil
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sigil-daemon
ExecStartPost=/bin/sh -c 'sleep 0.5 && chmod 666 /tmp/sigil.sock'
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=false
PrivateDevices=false
DeviceAllow=/dev/sd*
ReadWritePaths=/tmp

# Allow access to USB devices
SupplementaryGroups=sigil plugdev

# State directory
StateDirectory=sigil
ConfigurationDirectory=sigil

# Environment
Environment=SIGIL_CONFIG=/etc/sigil/daemon.json
Environment=RUST_LOG=sigil_daemon=info

[Install]
WantedBy=multi-user.target
